# resources/fpfis-php71-dev-7-ssh/resources/supervisor-sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -e
environment=DOCUMENT_ROOT=%(ENV_DOCUMENT_ROOT)s
